# ops/systemd/zutomayo-retention.timer
[Unit]
Description=Run the Zutomayo retention worker daily

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=zutomayo-retention.service

[Install]
WantedBy=timers.target
